# linux/config_backup/acer_vswitch/etc/rc.local
#!/bin/bash

# created with: 
# * fallocate -l 4G /swap.file
# * chmod 600 /swap.file
# * mkswap /swap.file 
swapon /swap.file
